# supervisor-climate-queue.conf
[program:laravel_climate_queue]
command=php /home/luizescobal/study/teste-dev-backend/artisan queue:work --queue=climate_data --sleep=3 --tries=3 --max-time=3600
directory=/home/luizescobal/study/teste-dev-backend
autostart=true
autorestart=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel_climate_queue.log
stderr_logfile=/var/log/supervisor/laravel_climate_queue_error.log
